# templates/github-hook-relay.ini
[program:github-hook-relay]
command: /opt/github-hook-relay/github-hook-relay
environment: JENKINS_ENDPOINT={{ relay_jenkins_endpoint }},AWS_SQS_QUEUE={{ relay_sqs_queue }}
