# templates/github-hook-relay.ini
[program:github-hook-relay]
command: /opt/github-hook-relay/github-hook-relay
